# r2d2-camera-stream.service
[Unit]
Description=R2D2 Camera Stream Service (MJPEG)
After=network.target
Documentation=file:///home/severin/dev/r2d2/001_ARCHITECTURE_OVERVIEW.md

[Service]
Type=simple
User=severin
Group=severin
WorkingDirectory=/home/severin/dev/r2d2/ros2_ws
Environment="PATH=/home/severin/depthai_env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="OPENBLAS_CORETYPE=ARMV8"
Environment="ROS_DOMAIN_ID=0"

# Use launch file (starts camera + stream together)
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source install/setup.bash && ros2 launch r2d2_camera camera_stream.launch.py'

# Auto-restart on failure (but wait longer to avoid rapid restarts)
Restart=on-failure
RestartSec=10s
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=r2d2-camera-stream

# Service hardening
NoNewPrivileges=true
ProtectHome=no
ReadWritePaths=/home/severin/dev/r2d2

[Install]
WantedBy=multi-user.target

