# r2d2-camera-perception.service
[Unit]
Description=R2D2 Camera Perception Service (Face Recognition)
After=network.target audio-notification.service
Documentation=file:///home/severin/dev/r2d2/QUICK_START.md

[Service]
Type=exec
User=severin
WorkingDirectory=/home/severin/dev/r2d2/ros2_ws
Environment="PATH=/home/severin/depthai_env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Source ROS 2 setup and launch camera
ExecStart=/bin/bash -c 'source install/setup.bash && ros2 launch r2d2_bringup r2d2_camera_perception.launch.py enable_face_recognition:=true enable_gesture_recognition:=true gesture_recognition_model_path:=/home/severin/dev/r2d2/data/gesture_recognition/models/severin_gesture_classifier.pkl'

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=r2d2-camera

[Install]
WantedBy=multi-user.target
