# ros2_ws/src/r2d2_audio/config/r2d2-volume-control.service
[Unit]
Description=R2D2 Volume Control Node - Master volume publisher
Documentation=file:///home/severin/dev/r2d2/260_VOLUME_CONTROL_REFERENCE.md
After=network.target

# Start before other audio services so they receive initial volume
Before=r2d2-audio-notification.service r2d2-gesture-intent.service r2d2-speech-node.service

[Service]
Type=simple
User=severin
Group=severin
WorkingDirectory=/home/severin/dev/r2d2

# Environment setup
Environment="HOME=/home/severin"
Environment="ROS_DOMAIN_ID=0"
Environment="PYTHONUNBUFFERED=1"

# Source ROS2 and workspace, then run volume control node
# Software-only mode: hardware_enabled=False (no ADC required)
ExecStart=/bin/bash -c '\
    source /opt/ros/humble/setup.bash && \
    source /home/severin/dev/r2d2/ros2_ws/install/setup.bash && \
    ros2 run r2d2_audio volume_control_node --ros-args \
        -p hardware_enabled:=false \
        --params-file /home/severin/dev/r2d2/ros2_ws/src/r2d2_audio/config/audio_params.yaml'

# Restart policy
Restart=always
RestartSec=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=r2d2-volume-control

[Install]
WantedBy=multi-user.target
